# src/main/assembly/dist/install/easy-deposit-report.service
[Unit]
Description=EASY Deposit Report Service

[Service]
ExecStart=/bin/java \
   -Dlogback.configurationFile=/etc/opt/dans.knaw.nl/easy-deposit-report/logback-service.xml \
   -Dapp.home=/opt/dans.knaw.nl/easy-deposit-report \
   -jar /opt/dans.knaw.nl/easy-deposit-report/bin/easy-deposit-report.jar run-service

User=easy-deposit-report
Group=easy-deposit-report

[Install]
WantedBy=multi-user.target
